# etc/systemd/system/mounter.service
[Unit]
Description=Mount required stuff

[Service]
ExecStart=/usr/local/bin/mounts

[Install]
WantedBy=local-fs.target
